# archisoarm_releng/aarch64/airootfs/var/lib/pacman/local/pcsclite-2.0.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/pcsc-spy
usr/bin/pcscd
usr/include/
usr/include/PCSC/
usr/include/PCSC/debuglog.h
usr/include/PCSC/ifdhandler.h
usr/include/PCSC/pcsclite.h
usr/include/PCSC/reader.h
usr/include/PCSC/winscard.h
usr/include/PCSC/wintypes.h
usr/lib/
usr/lib/libpcsclite.so
usr/lib/libpcsclite.so.1
usr/lib/libpcsclite.so.1.0.0
usr/lib/libpcscspy.so
usr/lib/libpcscspy.so.0
usr/lib/libpcscspy.so.0.0.0
usr/lib/pcsc/
usr/lib/pcsc/drivers/
usr/lib/pkgconfig/
usr/lib/pkgconfig/libpcsclite.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/pcscd.service
usr/lib/systemd/system/pcscd.socket
usr/share/
usr/share/doc/
usr/share/doc/pcsc-lite/
usr/share/doc/pcsc-lite/README.polkit
usr/share/doc/pcsc-lite/install_spy.sh
usr/share/doc/pcsc-lite/uninstall_spy.sh
usr/share/licenses/
usr/share/licenses/pcsclite/
usr/share/licenses/pcsclite/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/pcsc-spy.1.gz
usr/share/man/man5/
usr/share/man/man5/reader.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/pcscd.8.gz
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.debian.pcsc-lite.policy

