# archisoarm_releng/aarch64/airootfs/var/lib/pacman/local/rpcbind-1.2.6-3/files
%FILES%
etc/
etc/conf.d/
etc/conf.d/rpcbind
usr/
usr/bin/
usr/bin/rpcbind
usr/bin/rpcinfo
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rpcbind.service
usr/lib/systemd/system/rpcbind.socket
usr/lib/sysusers.d/
usr/lib/sysusers.d/rpcbind.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/rpcbind.conf
usr/share/
usr/share/licenses/
usr/share/licenses/rpcbind/
usr/share/licenses/rpcbind/COPYING
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/rpcbind.8.gz
usr/share/man/man8/rpcinfo.8.gz

%BACKUP%
etc/conf.d/rpcbind	eed7d755165ad49a5f7015a802f096a1

