# archisoarm_releng/aarch64/airootfs/var/lib/pacman/local/modemmanager-1.20.6-1/files
%FILES%
etc/
etc/ModemManager/
etc/ModemManager/connection.d/
etc/ModemManager/fcc-unlock.d/
usr/
usr/bin/
usr/bin/ModemManager
usr/bin/mmcli
usr/lib/
usr/lib/ModemManager/
usr/lib/ModemManager/connection.d/
usr/lib/ModemManager/fcc-unlock.d/
usr/lib/ModemManager/libmm-plugin-altair-lte.so
usr/lib/ModemManager/libmm-plugin-anydata.so
usr/lib/ModemManager/libmm-plugin-broadmobi.so
usr/lib/ModemManager/libmm-plugin-cinterion.so
usr/lib/ModemManager/libmm-plugin-dell.so
usr/lib/ModemManager/libmm-plugin-dlink.so
usr/lib/ModemManager/libmm-plugin-ericsson-mbm.so
usr/lib/ModemManager/libmm-plugin-fibocom.so
usr/lib/ModemManager/libmm-plugin-foxconn.so
usr/lib/ModemManager/libmm-plugin-generic.so
usr/lib/ModemManager/libmm-plugin-gosuncn.so
usr/lib/ModemManager/libmm-plugin-haier.so
usr/lib/ModemManager/libmm-plugin-huawei.so
usr/lib/ModemManager/libmm-plugin-intel.so
usr/lib/ModemManager/libmm-plugin-iridium.so
usr/lib/ModemManager/libmm-plugin-linktop.so
usr/lib/ModemManager/libmm-plugin-longcheer.so
usr/lib/ModemManager/libmm-plugin-motorola.so
usr/lib/ModemManager/libmm-plugin-mtk.so
usr/lib/ModemManager/libmm-plugin-nokia-icera.so
usr/lib/ModemManager/libmm-plugin-nokia.so
usr/lib/ModemManager/libmm-plugin-novatel-lte.so
usr/lib/ModemManager/libmm-plugin-novatel.so
usr/lib/ModemManager/libmm-plugin-option-hso.so
usr/lib/ModemManager/libmm-plugin-option.so
usr/lib/ModemManager/libmm-plugin-pantech.so
usr/lib/ModemManager/libmm-plugin-quectel.so
usr/lib/ModemManager/libmm-plugin-samsung.so
usr/lib/ModemManager/libmm-plugin-sierra-legacy.so
usr/lib/ModemManager/libmm-plugin-sierra.so
usr/lib/ModemManager/libmm-plugin-simtech.so
usr/lib/ModemManager/libmm-plugin-telit.so
usr/lib/ModemManager/libmm-plugin-thuraya.so
usr/lib/ModemManager/libmm-plugin-tplink.so
usr/lib/ModemManager/libmm-plugin-ublox.so
usr/lib/ModemManager/libmm-plugin-via.so
usr/lib/ModemManager/libmm-plugin-wavecom.so
usr/lib/ModemManager/libmm-plugin-x22x.so
usr/lib/ModemManager/libmm-plugin-zte.so
usr/lib/ModemManager/libmm-shared-fibocom.so
usr/lib/ModemManager/libmm-shared-foxconn.so
usr/lib/ModemManager/libmm-shared-icera.so
usr/lib/ModemManager/libmm-shared-novatel.so
usr/lib/ModemManager/libmm-shared-option.so
usr/lib/ModemManager/libmm-shared-sierra.so
usr/lib/ModemManager/libmm-shared-telit.so
usr/lib/ModemManager/libmm-shared-xmm.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ModemManager.service
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/77-mm-broadmobi-port-types.rules
usr/lib/udev/rules.d/77-mm-cinterion-port-types.rules
usr/lib/udev/rules.d/77-mm-dell-port-types.rules
usr/lib/udev/rules.d/77-mm-dlink-port-types.rules
usr/lib/udev/rules.d/77-mm-ericsson-mbm.rules
usr/lib/udev/rules.d/77-mm-fibocom-port-types.rules
usr/lib/udev/rules.d/77-mm-foxconn-port-types.rules
usr/lib/udev/rules.d/77-mm-gosuncn-port-types.rules
usr/lib/udev/rules.d/77-mm-haier-port-types.rules
usr/lib/udev/rules.d/77-mm-huawei-net-port-types.rules
usr/lib/udev/rules.d/77-mm-linktop-port-types.rules
usr/lib/udev/rules.d/77-mm-longcheer-port-types.rules
usr/lib/udev/rules.d/77-mm-mtk-port-types.rules
usr/lib/udev/rules.d/77-mm-nokia-port-types.rules
usr/lib/udev/rules.d/77-mm-quectel-port-types.rules
usr/lib/udev/rules.d/77-mm-sierra.rules
usr/lib/udev/rules.d/77-mm-simtech-port-types.rules
usr/lib/udev/rules.d/77-mm-telit-port-types.rules
usr/lib/udev/rules.d/77-mm-tplink-port-types.rules
usr/lib/udev/rules.d/77-mm-ublox-port-types.rules
usr/lib/udev/rules.d/77-mm-x22x-port-types.rules
usr/lib/udev/rules.d/77-mm-zte-port-types.rules
usr/lib/udev/rules.d/80-mm-candidate.rules
usr/share/
usr/share/ModemManager/
usr/share/ModemManager/connection.available.d/
usr/share/ModemManager/connection.available.d/99-log-event
usr/share/ModemManager/fcc-unlock.available.d/
usr/share/ModemManager/fcc-unlock.available.d/03f0:4e1d
usr/share/ModemManager/fcc-unlock.available.d/105b
usr/share/ModemManager/fcc-unlock.available.d/105b:e0ab
usr/share/ModemManager/fcc-unlock.available.d/1199
usr/share/ModemManager/fcc-unlock.available.d/1199:9079
usr/share/ModemManager/fcc-unlock.available.d/1eac
usr/share/ModemManager/fcc-unlock.available.d/1eac:1001
usr/share/ModemManager/fcc-unlock.available.d/2c7c
usr/share/ModemManager/fcc-unlock.available.d/2c7c:030a
usr/share/ModemManager/fcc-unlock.available.d/413c:81a3
usr/share/ModemManager/fcc-unlock.available.d/413c:81a8
usr/share/ModemManager/mm-foxconn-t77w968-carrier-mapping.conf
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/mmcli
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Bearer.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Call.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Firmware.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Location.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Messaging.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Modem3gpp.ProfileManager.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Modem3gpp.Ussd.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Modem3gpp.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.ModemCdma.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Oma.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Sar.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Signal.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Simple.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Time.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.Voice.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Modem.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Sim.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.Sms.xml
usr/share/dbus-1/interfaces/org.freedesktop.ModemManager1.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.ModemManager1.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.freedesktop.ModemManager1.conf
usr/share/gtk-doc/
usr/share/gtk-doc/html/
usr/share/gtk-doc/html/ModemManager/
usr/share/gtk-doc/html/ModemManager/ModemManager-API-break-replacements.html
usr/share/gtk-doc/html/ModemManager/ModemManager-Common-udev-tags.html
usr/share/gtk-doc/html/ModemManager/ModemManager-Errors.html
usr/share/gtk-doc/html/ModemManager/ModemManager-Flags-and-Enumerations.html
usr/share/gtk-doc/html/ModemManager/ModemManager-Version-checks.html
usr/share/gtk-doc/html/ModemManager/ModemManager-interface-initialization-sequence-subclassed.png
usr/share/gtk-doc/html/ModemManager/ModemManager-interface-initialization-sequence.png
usr/share/gtk-doc/html/ModemManager/ModemManager-logo-square.png
usr/share/gtk-doc/html/ModemManager/ModemManager-logo-wide-text.png
usr/share/gtk-doc/html/ModemManager/ModemManager-logo-wide.png
usr/share/gtk-doc/html/ModemManager/ModemManager-states.png
usr/share/gtk-doc/html/ModemManager/ModemManager.devhelp2
usr/share/gtk-doc/html/ModemManager/api-index-1-0.html
usr/share/gtk-doc/html/ModemManager/api-index-1-10.html
usr/share/gtk-doc/html/ModemManager/api-index-1-12.html
usr/share/gtk-doc/html/ModemManager/api-index-1-14.html
usr/share/gtk-doc/html/ModemManager/api-index-1-16.html
usr/share/gtk-doc/html/ModemManager/api-index-1-18.html
usr/share/gtk-doc/html/ModemManager/api-index-1-2.html
usr/share/gtk-doc/html/ModemManager/api-index-1-20.html
usr/share/gtk-doc/html/ModemManager/api-index-1-4.html
usr/share/gtk-doc/html/ModemManager/api-index-1-6.html
usr/share/gtk-doc/html/ModemManager/api-index-1-8.html
usr/share/gtk-doc/html/ModemManager/api-index-deprecated.html
usr/share/gtk-doc/html/ModemManager/api-index-full.html
usr/share/gtk-doc/html/ModemManager/ch02s02.html
usr/share/gtk-doc/html/ModemManager/ch03s02.html
usr/share/gtk-doc/html/ModemManager/ch04s02.html
usr/share/gtk-doc/html/ModemManager/ch04s03.html
usr/share/gtk-doc/html/ModemManager/ch04s04.html
usr/share/gtk-doc/html/ModemManager/ch06s02.html
usr/share/gtk-doc/html/ModemManager/ch06s03.html
usr/share/gtk-doc/html/ModemManager/ch06s04.html
usr/share/gtk-doc/html/ModemManager/ch16.html
usr/share/gtk-doc/html/ModemManager/ch16s02.html
usr/share/gtk-doc/html/ModemManager/ch16s03.html
usr/share/gtk-doc/html/ModemManager/ch16s04.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Bearer.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Call.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Firmware.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Location.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Messaging.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Modem3gpp.ProfileManager.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Modem3gpp.Ussd.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Modem3gpp.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.ModemCdma.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Oma.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Sar.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Signal.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Simple.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Time.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.Voice.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Modem.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Sim.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.Sms.html
usr/share/gtk-doc/html/ModemManager/gdbus-org.freedesktop.ModemManager1.html
usr/share/gtk-doc/html/ModemManager/home.png
usr/share/gtk-doc/html/ModemManager/index.html
usr/share/gtk-doc/html/ModemManager/left-insensitive.png
usr/share/gtk-doc/html/ModemManager/left.png
usr/share/gtk-doc/html/ModemManager/ref-common-types.html
usr/share/gtk-doc/html/ModemManager/ref-compat.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-bus-name.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-object-bearer.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-object-call.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-object-manager.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-object-modem.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-object-sim.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-object-sms.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-standard-interfaces-introspectables.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-standard-interfaces-objectmanager.html
usr/share/gtk-doc/html/ModemManager/ref-dbus-standard-interfaces.html
usr/share/gtk-doc/html/ModemManager/ref-dbus.html
usr/share/gtk-doc/html/ModemManager/ref-migrating.html
usr/share/gtk-doc/html/ModemManager/ref-overview-introduction.html
usr/share/gtk-doc/html/ModemManager/ref-overview-modem-detection.html
usr/share/gtk-doc/html/ModemManager/ref-overview-modem-filter.html
usr/share/gtk-doc/html/ModemManager/ref-overview-modem-object-creation.html
usr/share/gtk-doc/html/ModemManager/ref-overview-modem-port-probing.html
usr/share/gtk-doc/html/ModemManager/ref-overview-modem-state-machine.html
usr/share/gtk-doc/html/ModemManager/ref-overview-plugin-specific-modems.html
usr/share/gtk-doc/html/ModemManager/ref-overview.html
usr/share/gtk-doc/html/ModemManager/ref-udev.html
usr/share/gtk-doc/html/ModemManager/right-insensitive.png
usr/share/gtk-doc/html/ModemManager/right.png
usr/share/gtk-doc/html/ModemManager/style.css
usr/share/gtk-doc/html/ModemManager/up-insensitive.png
usr/share/gtk-doc/html/ModemManager/up.png
usr/share/icons/
usr/share/icons/hicolor/
usr/share/icons/hicolor/22x22/
usr/share/icons/hicolor/22x22/apps/
usr/share/icons/hicolor/22x22/apps/ModemManager.png
usr/share/locale/
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/ModemManager.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/ModemManager.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/ModemManager.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/ModemManager.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/ModemManager.mo
usr/share/locale/fur/
usr/share/locale/fur/LC_MESSAGES/
usr/share/locale/fur/LC_MESSAGES/ModemManager.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/ModemManager.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/ModemManager.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/ModemManager.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/ModemManager.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/ModemManager.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/ModemManager.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/ModemManager.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/ModemManager.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/ModemManager.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/ModemManager.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/ModemManager.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/ModemManager.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/ModemManager.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/ModemManager.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/ModemManager.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/ModemManager.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/mmcli.1.gz
usr/share/man/man8/
usr/share/man/man8/ModemManager.8.gz
usr/share/polkit-1/
usr/share/polkit-1/actions/
usr/share/polkit-1/actions/org.freedesktop.ModemManager1.policy

