# archisoarm_releng/aarch64/airootfs/var/lib/pacman/local/usbmuxd-1.1.1-3/files
%FILES%
usr/
usr/bin/
usr/bin/usbmuxd
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/usbmuxd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/usbmuxd.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/39-usbmuxd.rules
usr/share/
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/usbmuxd.8.gz

