# archisoarm_releng/aarch64/airootfs/var/lib/pacman/local/openvpn-2.6.5-1/files
%FILES%
usr/
usr/bin/
usr/bin/openvpn
usr/include/
usr/include/openvpn-msg.h
usr/include/openvpn-plugin.h
usr/lib/
usr/lib/openvpn/
usr/lib/openvpn/plugins/
usr/lib/openvpn/plugins/openvpn-plugin-auth-pam.so
usr/lib/openvpn/plugins/openvpn-plugin-down-root.so
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/openvpn-client@.service
usr/lib/systemd/system/openvpn-server@.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/openvpn.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/openvpn.conf
usr/share/
usr/share/doc/
usr/share/doc/openvpn/
usr/share/doc/openvpn/COPYING
usr/share/doc/openvpn/COPYRIGHT.GPL
usr/share/doc/openvpn/Changes.rst
usr/share/doc/openvpn/README
usr/share/doc/openvpn/README.auth-pam
usr/share/doc/openvpn/README.down-root
usr/share/doc/openvpn/README.mbedtls
usr/share/doc/openvpn/README.systemd
usr/share/doc/openvpn/gui-notes.txt
usr/share/doc/openvpn/management-notes.txt
usr/share/doc/openvpn/openvpn-examples.5.html
usr/share/doc/openvpn/openvpn.8.html
usr/share/licenses/
usr/share/licenses/openvpn/
usr/share/licenses/openvpn/COPYING
usr/share/licenses/openvpn/COPYRIGHT.GPL
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/openvpn-examples.5.gz
usr/share/man/man8/
usr/share/man/man8/openvpn.8.gz
usr/share/openvpn/
usr/share/openvpn/contrib/
usr/share/openvpn/contrib/OCSP_check/
usr/share/openvpn/contrib/OCSP_check/OCSP_check.sh
usr/share/openvpn/contrib/README
usr/share/openvpn/contrib/extract-crl/
usr/share/openvpn/contrib/extract-crl/extractcrl.py
usr/share/openvpn/contrib/multilevel-init.patch
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/README
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/fwmarkroute.down
usr/share/openvpn/contrib/openvpn-fwmarkroute-1.00/fwmarkroute.up
usr/share/openvpn/contrib/pull-resolv-conf/
usr/share/openvpn/contrib/pull-resolv-conf/client.down
usr/share/openvpn/contrib/pull-resolv-conf/client.up
usr/share/openvpn/contrib/vcpkg-ports/
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/0001-nmake-compatibility-with-vcpkg-nmake.patch
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/0002-config-w32-vc.h.in-indicate-OpenSSL.patch
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/CONTROL
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/pkcs11-helper-001-RFC7512.patch
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/pkcs11-helper-002-dynamic_loader_flags.patch
usr/share/openvpn/contrib/vcpkg-ports/pkcs11-helper/portfile.cmake
usr/share/openvpn/contrib/vcpkg-triplets/
usr/share/openvpn/contrib/vcpkg-triplets/arm64-windows-ovpn.cmake
usr/share/openvpn/contrib/vcpkg-triplets/x64-windows-ovpn.cmake
usr/share/openvpn/contrib/vcpkg-triplets/x86-windows-ovpn.cmake
usr/share/openvpn/examples/
usr/share/openvpn/examples/README
usr/share/openvpn/examples/client.conf
usr/share/openvpn/examples/firewall.sh
usr/share/openvpn/examples/home.up
usr/share/openvpn/examples/loopback-client
usr/share/openvpn/examples/loopback-server
usr/share/openvpn/examples/office.up
usr/share/openvpn/examples/openvpn-shutdown.sh
usr/share/openvpn/examples/openvpn-startup.sh
usr/share/openvpn/examples/server.conf
usr/share/openvpn/examples/tls-home.conf
usr/share/openvpn/examples/tls-office.conf

