# archisoarm_releng/aarch64/airootfs/usr/lib/systemd/user/gssuserproxy.service
[Unit]
Description=GSS User Proxy
Documentation=man:gssproxy(8)

[Service]
Type=exec
StandardError=journal
ExecStart=/usr/bin/gssproxy -i -u
Restart=on-failure
