# archisoarm_releng/aarch64/airootfs/var/lib/pacman/local/archlinux-keyring-20230704-1/files
%FILES%
usr/
usr/bin/
usr/bin/archlinux-keyring-wkd-sync
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/archlinux-keyring-wkd-sync.service
usr/lib/systemd/system/archlinux-keyring-wkd-sync.timer
usr/lib/systemd/system/timers.target.wants/
usr/lib/systemd/system/timers.target.wants/archlinux-keyring-wkd-sync.timer
usr/share/
usr/share/pacman/
usr/share/pacman/keyrings/
usr/share/pacman/keyrings/archlinux-revoked
usr/share/pacman/keyrings/archlinux-trusted
usr/share/pacman/keyrings/archlinux.gpg

